# nightly_task/crontab
0 3 * * * /usr/local/bin/python /app/script.py >> /app/logs/cron.log 2>&1
