# nightly_task/crontab
0 3 * * * python /app/script.py >> /app/logs/cron.log 2>&1
